# cron/crontab
*/2 * * * * python3 /credit_update.py >> /logfile
